# startup.conf
start on startup

exec forever start /vagrant/www/server.js