# bootstrapping/nodeserver.conf
#!upstart
description "node.js server process"
author      "Nate Marquardt"

start on startup
stop on shutdown

script
    export HOME="/root"

    echo $$ > /var/run/nodeserver.pid
    exec sudo -u root /usr/bin/supervisor -w /vagrant/www/routes,/vagrant/www/server.js /vagrant/www/server.js >> /var/log/nodeserver.sys.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/nodeserver.sys.log
end script

pre-stop script
    rm /var/run/nodeserver.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/nodeserver.sys.log
end script
